# packaging/systemd/tnv-full.service
[Unit]
Description=Tax N Vote — full (email/password) site
Documentation=https://github.com/greenpdx/TaxNVote26
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tnv
Group=tnv
WorkingDirectory=/opt/tnv/full
EnvironmentFile=/etc/tnv/full.env
ExecStart=/opt/tnv/full/tnv-server
Restart=on-failure
RestartSec=5

# Hardening — the app only needs to write to /var/lib/tnv (e.g. SQLite).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/tnv

[Install]
WantedBy=multi-user.target
